# scripts/nanoclaw-v2.service
[Unit]
Description=NanoClaw v2 Personal Assistant
After=network.target

[Service]
User=ubuntu
Group=ubuntu
Type=simple
ExecStart=/usr/bin/node /home/ubuntu/nanoclaw-v2/dist/index.js
# Belt-and-suspenders: the Node shutdown handler stops in-flight containers,
# but if Node crashes early or is killed by SIGKILL, ExecStop sweeps any
# v2 containers that escaped. Parallel xargs caps wall time; systemd runs
# ExecStop synchronously and TimeoutStopSec applies.
#
# Filter on nanoclaw-v2- prefix so we don't touch v1 containers while both
# services are running side-by-side during the migration window.
ExecStop=/bin/sh -c '/usr/bin/docker ps --filter name=nanoclaw-v2- --format "{{.Names}}" | xargs -r -P 10 -I{} /usr/bin/docker stop -t 10 {}'
TimeoutStopSec=60
WorkingDirectory=/home/ubuntu/nanoclaw-v2
Restart=always
RestartSec=5
Environment=HOME=/home/ubuntu
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/ubuntu/.local/bin
StandardOutput=append:/home/ubuntu/nanoclaw-v2/logs/nanoclaw.log
StandardError=append:/home/ubuntu/nanoclaw-v2/logs/nanoclaw.error.log

[Install]
WantedBy=multi-user.target
